# rustboard.service
[Unit]
Description=rustboard
After=network.target

[Service]
Type=simple
User=rustboard
WorkingDirectory=/opt/rustboard
ExecStart=/opt/rustboard/rustboard
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
